# deployment/chatwoot.target
[Unit]
Wants=unixp-web.1.service unixp-worker.1.service

[Install]
WantedBy=multi-user.target