# deployment/chatwoot-worker.target
[Unit]
Description=UniXP Background Worker
Wants=unixp-worker.1.service

[Install]
WantedBy=multi-user.target